# templates/systemd/rei-orchestrator-worker.timer
[Unit]
Description=Run Rei Orchestrator Worker every {{INTERVAL_SECONDS}}s

[Timer]
OnBootSec=30s
OnUnitActiveSec={{INTERVAL_SECONDS}}s
AccuracySec=1s
Unit=rei-orchestrator-worker.service

[Install]
WantedBy=timers.target
