# templates/systemd/rei-orchestrator-worker.service
[Unit]
Description=Rei Orchestrator Worker Tick
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User={{APP_USER}}
WorkingDirectory={{APP_DIR}}
Environment=OPENCLAW_HOME={{OPENCLAW_HOME}}
Environment=PYTHONPATH={{APP_DIR}}
ExecStart=/usr/bin/env bash -lc 'cd {{APP_DIR}} && python3 -m src.orchestrator.cli worker-tick'
